# .config/autostart/wsgidav.desktop
[Desktop Entry]
Exec=/usr/bin/wsgidav -c ~/wsgidav.plain.yaml
Icon=folder-publicshare
Name=wsgidav
Terminal=False
Type=Application
